# azure-ollama/rootfs/etc/services.d/azure-ollama/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Azure Ollama
# Executes the Node.js application.
# ==============================================================================
bashio::log.info "Starting Azure Ollama Node.js application..."

serverPort=$(bashio::addon.port 8099)
bashio::log.info "server port ${serverPort}"
export SERVER_PORT=${serverPort:=8099}

node /app/server.js
