# local_google/rootfs/etc/services.d/google/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Google Proxy Service..."

serverPort=$(bashio::addon.port)

udpPort=$(bashio::addon.port 3312)

ingressUrl="$(bashio::addon.ingress_url)"

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "server port ${serverPort}"
bashio::log.info "udp port ${udpPort}"
bashio::log.info "ingressUrl ${ingressUrl}"

export SERVER_PORT=${serverPort:=8089}
export UDP_PORT=${udpPort:=3312}
export INGRESS_URL=${ingressUrl:=""}

cd /www

exec node server/index.js
