# local_google/rootfs/etc/services.d/google/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Google Proxy Service..."

cd /www

exec node server/index.js
